# debian/smartass-daemon.service
[Unit]
Description=Smartass daemon (plugin host + D-Bus service)
After=graphical-session.target

[Service]
Type=dbus
BusName=ai.talonic.Smartass
ExecStart=/opt/smartass/bin/python -m smartass.daemon
Restart=on-failure
RestartSec=3
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
